# .config/systemd/user/blueman-applet.service
[Service]
ExecStart=/usr/bin/blueman-applet
Restart=always
RestartSec=1

[Unit]
StopWhenUnneeded=true
