# .config/systemd/user/desktop.target
[Unit]
Wants=light-locker.service
Wants=polybar.service
Wants=xfce4-volumed.service
