# .config/systemd/user/desktop.target
[Install]
Wants=light-locker.service
Wants=polybar.service
