# .config/systemd/user/mate-power-manager.service
[Service]
Type=simple
ExecStart=/usr/bin/mate-power-manager
Restart=always
RestartSec=1

[Unit]
StopWhenUnneeded=true
