# .config/systemd/user/nm-applet.service
[Service]
ExecStart=/usr/bin/nm-applet
Restart=always
RestartSec=1

[Unit]
StopWhenUnneeded=true