# .config/systemd/user/system-tray.target
[Unit]
Wants=blueman-applet.service
Wants=dropbox.service
Wants=nm-applet.service
Wants=xfce4-power-manager.service