# .config/systemd/user/xfce4-volumed.service
[Service]
Type=simple
ExecStart=/usr/bin/xfce4-volumed --no-daemon
Restart=always
RestartSec=1

[Unit]
StopWhenUnneeded=true
